# keeper.service
[Unit]
Description=Keeper - note-taking app
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/src/keeper
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/usr/bin/env npm run build
ExecStart=/usr/bin/env npm run server
Restart=on-failure
RestartSec=5
Environment=PORT=3001
Environment=DATA_DIR=%h/src/keeper/data

[Install]
WantedBy=default.target
