# keeper.service
[Unit]
Description=Keeper - note-taking app
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/jackkendall/src/keeper
Environment=PATH=/home/jackkendall/.local/share/nvm/v24.15.0/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/home/jackkendall/.local/share/nvm/v24.15.0/bin/npm run build
ExecStart=/home/jackkendall/.local/share/nvm/v24.15.0/bin/npx tsx server/index.ts
Restart=on-failure
RestartSec=5
Environment=PORT=3001
Environment=DATA_DIR=/home/jackkendall/src/keeper/data

[Install]
WantedBy=default.target
